# docker/airflow/runit/worker/run
#!/bin/sh -e
exec 2>&1

if [ "$(id -u)" = "0" ]; then
    exec chpst -u airflow airflow worker
else
    exec airflow worker
fi
