# docker/airflow/runit/worker/run
#!/bin/sh -e
exec 2>&1

if [ "$(id -u)" = "0" ]; then
    exec chpst -u airflow airflow celery worker
else
    exec airflow celery worker
fi
